# phases/phase-1-fedtracker-migration/ansible/templates/fedtracker.service.j2
[Unit]
Description=FedTracker Personnel Tracking Application (FastAPI)
After=network.target

[Service]
User={{ app_user }}
Group={{ app_user }}
WorkingDirectory={{ app_dir }}
Environment=DB_TYPE=sqlite
Environment=SQLITE_PATH={{ app_dir }}/fedtracker.db
ExecStart=/usr/bin/python3 -m uvicorn main:app --host 0.0.0.0 --port {{ app_port }}
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security hardening (same directives from Day 1, now automated)
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths={{ app_dir }}
PrivateTmp=yes
MemoryMax=512M
OOMScoreAdjust=-500

[Install]
WantedBy=multi-user.target